# local.llama-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — do not install this file directly.

  Run ./setup-tunnel.sh instead. It will:
    1. Prompt for SSH host and port forwarding spec via macOS dialogs
    2. Store config in macOS Keychain (no secrets in files)
    3. Install launch-tunnel.sh to ~/.local/bin/
    4. Generate and load the real plist at ~/Library/LaunchAgents/

  To uninstall: ./uninstall-tunnel.sh

  This template shows the structure for reference only.
  The generated plist invokes launch-tunnel.sh, which reads
  SSH host and port config from Keychain at runtime.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.llama-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Generated by setup-tunnel.sh with absolute path to wrapper -->
    <string>~/.local/bin/launch-llama-tunnel.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/llama-tunnel.err</string>
  <key>StandardOutPath</key>
  <string>/tmp/llama-tunnel.out</string>
</dict>
</plist>
